# deploy/upstart.conf
description "aniairadek sites virtualenv and gunicorn_django"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /var/web/aniairadek.info/uwsgi_server.sh